# build/supervisor/messenger-low.conf
[program:messenger-low]
directory     = /var/www/html
command       = php bin/console messenger:consume async_priority_low -vv
numprocs      = 4
process_name  = %(program_name)s_%(process_num)02d
autostart     = true
autorestart   = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
