# build/supervisor/messenger-worker.conf
[supervisord]
nodaemon=true

[program:messenger-worker]
command=php /var/www/html/bin/console messenger:consume async_priority_low -vv
directory=/var/www/html
autostart=true
autorestart=true
startsecs=5

# send worker’s STDOUT into the container’s stdout
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

# send worker’s STDERR into the container’s stderr
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# optional: merge stderr into stdout
redirect_stderr=true
