# build/supervisor/messenger-worker.conf
[program:messenger-worker]
directory=/var/www/html
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
command=php /var/www/html/bin/console messenger:consume async_concurrent --time-limit=3600 --memory-limit=128M --no-interaction -vv
numprocs=4
process_name=%(program_name)s_%(process_num)02d
