# build/supervisor/messenger-concurrent.conf
[program:messenger-concurrent]
directory     = /var/www/html
# command       = php /var/www/html/bin/console messenger:consume async_concurrent --sleep=100000 -vv
command       = php /var/www/html/bin/console messenger:consume async_concurrent -vv
numprocs      = 6                # many tiny jobs can run in parallel
process_name  = %(program_name)s_%(process_num)02d
autostart     = true
autorestart   = true
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
